# rootfs/etc/services.d/simplelogin-email/run
#!/command/with-contenv bash
# shellcheck shell=bash
set -euo pipefail
source /etc/simplelogin-aio/env-helpers.sh
normalize_simplelogin_blank_sensitive_env_vars
normalize_simplelogin_compat_env_vars
normalize_simplelogin_presence_env_vars

until curl -fsS http://127.0.0.1:7777/health >/dev/null 2>&1; do
    echo "Waiting for SimpleLogin web service before starting email handler..."
    sleep 2
done

cd /code || exit 1
exec_as_simplelogin_with_sanitized_env python email_handler.py
