# rootfs/etc/services.d/simplelogin-job/run
#!/command/with-contenv bash
cd /code || exit 1
exec s6-setuidgid simplelogin python job_runner.py